# cam_stack/rootfs/etc/services.d/go2rtc/run
#!/usr/bin/with-contenv bashio
bashio::log.info "starting go2rtc"
exec go2rtc -config /etc/go2rtc.yaml
